# etc/cron.daily/logrotate
#!/bin/sh

logrotate --state="$XDG_STATE_HOME/logrotate/logrotate.status" "$HOME/etc/logrotate.conf"

