# etc/cron.daily/logrotate
#!/bin/sh

logrotate --state="$HOME/var/lib/logrotate/logrotate.status" "$HOME/etc/logrotate.conf"

