# etc/cron.daily/updatedb
#!/bin/sh

updatedb -l 0 -n ".git" -U "$HOME" -o "$HOME/var/lib/locate/locate.db"

