# XDG_CONFIG_HOME/cron.daily/updatedb
#!/bin/sh

updatedb -l0 -n .git -U "$HOME" -o "$LOCATEDB"

