# config/systemd/user/atworkd.service
[Unit]
Description=atworkd
Wants=network-online.target

[Service]
ExecStart=%h/bin/atworkd

[Install]
WantedBy=multi-user.target
